# module/bluetooth/cpc-bluetooth/debconf/services/cpc-bluetooth.service.in
[Unit]
Description=CPC HCI Bridge service (ver @PROJECT_VER@)


[Service]
WorkingDirectory=/var/lib/rafael/dev
ExecStartPre=${CPC_BLUETOOTH_SERVICE_PATH}/stop
ExecStart=${CPC_BLUETOOTH_SERVICE_PATH}/cpc-bluetooth
ExecStartPost=${CPC_BLUETOOTH_SERVICE_PATH}/startpost
ExecStopPost=${CPC_BLUETOOTH_SERVICE_PATH}/stop
# Restart=on-failure
KillMode=process
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
