# integrate/debian/border_router/debconf/services/ez-otbr.service.in
[Unit]
Description=EZMESH Openthread Border Router (ver @PROJECT_VER@)

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/ez-otbr
ExecStart=${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_DATAROOTDIR}/otbr/script/server

[Install]
WantedBy=multi-user.target

