# module/cpc/cpcd-system/debconf/services/cpc-mgmt.service.in
[Unit]
Description=CPC Mgmt Daemon (ver @PROJECT_VER@)
After=bluetooth.service otbr-agent.service dbus.service network.target systemd-timesyncd.service hciuart.service serial-getty.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=${CMAKE_INSTALL_PREFIX}/src/scripts/cpc-mgmt 
ExecReload=/bin/kill -HUP $MAINPID
# Restart=on-failure
# RestartSec=5s

[Install]
WantedBy=multi-user.target
