# module/cpc/cpcd-system/debconf/services/cpcd.service.in
[Unit]
Description=CPC Daemon (ver @PROJECT_VER@)

[Service]
ExecStart=${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_BINDIR}/cpcd -c ${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_SYSCONFDIR}/cpcd.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
KillMode=process
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
