# module/cpc/cpcd-system/debconf/postinst.in
#!/bin/sh -e
# postinst script for rafael
#
# see: dh_installdeb(1)

. /usr/share/debconf/confmodule

RAFAEL_VAR_DIR="${CPACK_PACKAGING_INSTALL_PREFIX}/var/lib/rafael/"
RAFAEL_CONFFILE=$RAFAEL_CONFDIR/rafael.cfg
CPCD_CONFFILE=${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_SYSCONFDIR}/cpcd.conf

# Defaults for required arguments
CPCD_DEFAULT_CFG_SERIAL="/dev/ttyUSB0"

setup_user() {
  if ! getent group rafael >/dev/null; then
    addgroup --quiet --system rafael
  fi

  if ! getent passwd rafael >/dev/null; then
    adduser --quiet --system --no-create-home --ingroup rafael --home "$RAFAEL_VAR_DIR" --shell /usr/sbin/nologin rafael
    usermod -a -G dialout rafael
  fi
}

fix_permissions() {
  mkdir -p $RAFAEL_VAR_DIR
  chown rafael $RAFAEL_VAR_DIR
}

enable_start_cpcd() {
  if type "systemctl" >/dev/null; then
    systemctl --system daemon-reload || true
    # systemctl start cpcd || true
    # systemctl enable cpcd || true
    systemctl start cpc-mgmt.service || true
    systemctl enable cpc-mgmt.service
  fi
}

handle_conf() {
  if ! grep -q "cpcd:" $RAFAEL_CONFFILE; then
    mkdir -p $RAFAEL_CONFDIR
    echo "cpcd:" >> $RAFAEL_CONFFILE
    echo "  serial: $CPCD_DEFAULT_CFG_SERIAL" >> $RAFAEL_CONFFILE
  fi

  db_get cpcd/serial_port
  CFG_CPCD_SERIAL="$RET"
  if [ "$RET" = "" ]; then
    CFG_CPCD_SERIAL="$CPCD_DEFAULT_CFG_SERIAL"
  fi

  # Update configuration file with outcome of configuration
  # This is a small state machine that can update the Unify config file, which is written in YAML
  state=""

  while IFS= read line #Read IFS https://en.wikipedia.org/wiki/Input_Field_Separators
  do
    # Set the state depending on the line read from the config file
    case "$line" in
      "cpcd:" )   # State: cpcd
        state="cpcd"
        ;;
      " "* | "-"* )   # Lines starting with '-' or ' ' don't trigger state change.
        # Do nothing.
        ;;
      *) state="" ;;    # Set state to empty, Any unknown text, that doesn't begin with '-' or ' '
    esac
    leftside=$(echo "$line" | cut -f1 -d":")
    # STM to update entries under 'cpcd:'
    case $state in
      "cpcd" )
        case "$line" in
        "cpcd:" )
          echo "$line" ;;                   # cpcd should not be indented
        *"serial:"*)
          echo "$leftside: $CFG_CPCD_SERIAL"
          # Update CPCd configuration file
          sed -i "/uart_device_file/c\uart_device_file: $CFG_CPCD_SERIAL" $CPCD_CONFFILE
          ;;
        *)
          echo "$line" ;;                 # Anything inside cpcd we indent
        esac ;;
      *)
        echo "$line" ;;                     # Anything else we just echo the line
                        # Without IFS= , echo here would remove the leading spaces and ruin the yaml formatting
    esac
  done < "$RAFAEL_CONFFILE" > "$RAFAEL_CONFFILE.tmp"
  mv "$RAFAEL_CONFFILE.tmp" "$RAFAEL_CONFFILE"

  echo "==="
  echo "$RAFAEL_CONFFILE new contents:"
  echo "==="
  cat "$RAFAEL_CONFFILE"
  echo "==="


}

case "$1" in
configure)
  handle_conf
  setup_user
  fix_permissions
  ;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

\
  *)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
  ;;
esac

#DEBHELPER#

# Enable and start RAFAEL_CPCD daemon
enable_start_cpcd

exit 0
