# module/cpc/cpcd-system/debconf/prerm.in
#!/bin/sh
if type "systemctl" > /dev/null; then
    systemctl stop cpcd.service
    systemctl disable cpcd.service
fi
