# integrate/debian/controller/debconf/postinst.in
#!/bin/sh -e
# postinst script for rafael
#
# see: dh_installdeb(1)

. /usr/share/debconf/confmodule

RAFAEL_VAR_DIR="${CPACK_PACKAGING_INSTALL_PREFIX}/var/lib/rafael/"
RAFAEL_CONFDIR="${CPACK_PACKAGING_INSTALL_PREFIX}/etc/rafael"
RAFAEL_CONFFILE="${RAFAEL_CONFDIR}/rafael.cfg"
ezmesh_CONFFILE="${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_SYSCONFDIR}/ez_config.ini"

# Defaults for required arguments
ezmesh_DEFAULT_CFG_SERIAL="/dev/ttyUSB0"

setup_user() {
  if ! getent group rafael >/dev/null; then
    addgroup --quiet --system rafael
  fi

  if ! getent passwd rafael >/dev/null; then
    adduser --quiet --system --no-create-home --ingroup rafael --home "$RAFAEL_VAR_DIR" --shell /usr/sbin/nologin rafael
    usermod -a -G dialout rafael
  fi
}

fix_permissions() {
  mkdir -p $RAFAEL_VAR_DIR
  chown rafael $RAFAEL_VAR_DIR
}

enable_start_ezmesh() {
  if type "systemctl" >/dev/null; then
    systemctl --system daemon-reload || true
    systemctl start ez-mgmt.service || true
    systemctl enable ez-mgmt.service
  fi
}

handle_conf() {
  echo "ezmesh dir:" >> $RAFAEL_CONFDIR
  mkdir -p $RAFAEL_CONFDIR
  if ! grep -q "ezmesh:" $RAFAEL_CONFFILE; then
    echo "ezmesh dir:" >> $RAFAEL_CONFDIR
    mkdir -p $RAFAEL_CONFDIR
    echo "ezmesh:" >> $RAFAEL_CONFFILE
    echo "  serial: $ezmesh_DEFAULT_CFG_SERIAL" >> $RAFAEL_CONFFILE
  fi

  db_get ezmesh/serial_port
  CFG_ezmesh_SERIAL="$RET"
  if [ "$RET" = "" ]; then
    CFG_ezmesh_SERIAL="$ezmesh_DEFAULT_CFG_SERIAL"
  fi

  # Update configuration file with outcome of configuration
  # This is a small state machine that can update the Unify config file, which is written in YAML
  state=""

  while IFS= read line #Read IFS https://en.wikipedia.org/wiki/Input_Field_Separators
  do
    # Set the state depending on the line read from the config file
    case "$line" in
      "ezmesh:" )   # State: ezmesh
        state="ezmesh"
        ;;
      " "* | "-"* )   # Lines starting with '-' or ' ' don't trigger state change.
        # Do nothing.
        ;;
      *) state="" ;;    # Set state to empty, Any unknown text, that doesn't begin with '-' or ' '
    esac
    leftside=$(echo "$line" | cut -f1 -d":")
    # STM to update entries under 'ezmesh:'
    case $state in
      "ezmesh" )
        case "$line" in
        "ezmesh:" )
          echo "$line" ;;                   # ezmesh should not be indented
        *"serial:"*)
          echo "$leftside: $CFG_ezmesh_SERIAL"
          # Update ezmesh configuration file
          sed -i "/uart_device_file/c\uart_device_file: $CFG_ezmesh_SERIAL" $ezmesh_CONFFILE
          ;;
        *)
          echo "$line" ;;                 # Anything inside ezmesh we indent
        esac ;;
      *)
        echo "$line" ;;                     # Anything else we just echo the line
                        # Without IFS= , echo here would remove the leading spaces and ruin the yaml formatting
    esac
  done < "$RAFAEL_CONFFILE" > "$RAFAEL_CONFFILE.tmp"
  mv "$RAFAEL_CONFFILE.tmp" "$RAFAEL_CONFFILE"

  echo "==="
  echo "$RAFAEL_CONFFILE new contents:"
  echo "==="
  cat "$RAFAEL_CONFFILE"
  echo "==="


}

case "$1" in
configure)
  handle_conf
  setup_user
  fix_permissions
  ;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

\
  *)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
  ;;
esac

#DEBHELPER#
# Enable and start RAFAEL_ezmesh daemon
enable_start_ezmesh

exit 0
