# integrate/debian/border_router/debconf/services/otbr-agent.service.d/override.conf
[Unit]
PartOf=ez-otbr.service
BindsTo=ezmeshd.service
After=ezmeshd.service
