# module/border_router/cpc-otbr/debconf/services/otbr-agent.service.d/override.conf
[Unit]
PartOf=cpc-otbr.service
BindsTo=cpcd.service
After=cpcd.service
